# packaging/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable --now sms-doorman.service >/dev/null 2>&1 || true

exit 0

