# packaging/prerm
#!/usr/bin/env bash
set -euo pipefail

APP="sms-doorman"
systemctl stop "${APP}.service" || true
systemctl disable "${APP}.service" || true
systemctl daemon-reload || true

